# com.emailnotifications.flask.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.emailnotifications.flask</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/yonteru/Documents/Email_Notifications/venv/bin/python3</string>
        <string>/Users/yonteru/Documents/Email_Notifications/app.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/yonteru/Documents/Email_Notifications/flask_server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yonteru/Documents/Email_Notifications/flask_server_error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/yonteru/Documents/Email_Notifications</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>

